# arkdep-build.d/manjaro-gnome-ge-ogui/overlay/post_install/usr/lib/systemd/system/resize-fs.service
[Unit]
Description=Manjaro Resize-FS

[Service]
Type=simple
ExecStart=/usr/bin/resize-fs

[Install]
WantedBy=default.target
